# monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="XAUTHORITY_FILE_PATH", RUN+="/bin/bash MONITOR_SCRIPT_PATH"
